# meta-app-framework/recipes-core/applaunchd/applaunchd/agl-app@.service
[Unit]
Description=Sandboxed %I

[Service]
Type=simple
User=agl-driver
ExecStart=%i
Environment=XDG_RUNTIME_DIR=/run/user/1001/
