# meta-app-framework/recipes-core/applaunchd/applaunchd/applaunchd.service
[Unit]
Wants=network.target polkit.service
After=network.target polkit.service

[Service]
Type=notify
User=applaunchd
Environment=XDG_DATA_DIRS=/usr/share
ExecStart=/usr/bin/applaunchd
Restart=on-failure

[Install]
WantedBy=multi-user.target
